# etc/cron.hourly/movierequests
#!/bin/bash
#This script calls the movie request subscript.

/home/movierequests
