# etc/cron.hourly/transcode
#!/bin/sh
# script to start transcoding

/home/convert_files.sh
