# testing/etc/cron.hourly/syncgdrive
#!/bin/sh
# script to start gdrive sync hourly

/home/gdrive
